# files/etc/systemd/pkgcacheclean.service
[Unit]
Description=Clean package cache timer (/var/cache/pacman)

[Service]
Type=simple
ExecStart=/usr/bin/pkgcacheclean 3
