# roles/automount/files/mnt-nas-downloads.automount
[Unit]
Description=NFS share NAS:downloads
Requires=network-online.target

[Automount]
Where=/mnt/nas/downloads
TimeoutIdleSec=301

[Install]
WantedBy=multi-user.target
