# roles/automount/files/mnt-nas-backup.automount
[Unit]
Description=NFS share NAS:Backup
Requires=network-online.target

[Automount]
Where=/mnt/nas/backup
TimeoutIdleSec=301

[Install]
WantedBy=multi-user.target
