# roles/automount/files/mnt-nas-distrib.automount
[Unit]
Description=NFS share NAS:Distrib
Requires=network-online.target

[Automount]
Where=/mnt/nas/distrib
TimeoutIdleSec=301

[Install]
WantedBy=multi-user.target
