# roles/pkgcacheclean/files/pkgcacheclean.timer
[Unit]
Description=Clean package cache timer (/var/cache/pacman)

[Timer]
OnBootSec=30min
OnUnitActiveSec=1w 

[Install]
WantedBy=timers.target